# _fw-dump-dut.bin.extracted/squashfs-root/etc/uci-defaults/skworker
#!/bin/sh
chmod +x /etc/init.d/skworker
/etc/init.d/skworker enable
rm -rf /etc/uci-defaults/skworker
